# just-plane-mosher.service
[Unit]
Description=Just Plane Mosher Flight Tracker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=devinbernosky
WorkingDirectory=/home/devinbernosky/just-plane-mosher
ExecStart=/home/devinbernosky/just-plane-mosher/venv/bin/python -m src.main
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
